# supervisor/exhibitor.conf
[program:exhibitor]
command=/opt/exhibitor/wrapper.sh
autostart=true
autorestart=true
